# spdmd/service_files/spdmd.service
[Unit]
Description=Phosphor SPDM Daemon
After=obmc-mapper.target
After=mctp-ctrl.service
[Service]
Restart=always
Type=dbus
BusName=xyz.openbmc_project.SPDM
ExecStart=/usr/bin/spdmd

[Install]
WantedBy=multi-user.target
